# config/raspberry-pi/crontab
# /etc/cron.d/autohunter
# Cron jobs para manutenção do AutoHunter

# Limpeza operacional (logs/eventos/jobs/runs/notificações) a cada 6 horas
0 */6 * * * autohunter /home/autohunter/autohunter/venv/bin/python /home/autohunter/autohunter/scripts/cleanup_operational_data.py --apply

# Backup database (todos os dias 2am)
0 2 * * * autohunter /home/autohunter/autohunter/scripts/backup_db.sh

# Restart serviço (toda semana)
0 4 * * 0 root systemctl restart autohunter.service

# Cleanup de logs antigos (mensal)
0 0 1 * * autohunter find /home/autohunter/autohunter/logs -name "*.log.*" -mtime +30 -delete
