# config/raspberry-pi/crontab
# /etc/cron.d/autohunter
# Cron jobs para manutenńŃo do AutoHunter

# Limpeza de cache a cada 6 horas
0 */6 * * * autohunter /home/autohunter/autohunter/venv/bin/python -c "from scripts.cache_manager import get_cache; get_cache().cleanup_expired()"

# Vacuum database (toda segunda 3am)
0 3 * * 1 autohunter /home/autohunter/autohunter/venv/bin/python scripts/database_optimizer.py vacuum

# Backup database (todos os dias 2am)
0 2 * * * autohunter /home/autohunter/autohunter/scripts/backup_db.sh

# Restart servińo (toda semana)
0 4 * * 0 root systemctl restart autohunter.service

# Cleanup de logs antigos (mensal)
0 0 1 * * autohunter find /home/autohunter/autohunter/logs -name "*.log.*" -mtime +30 -delete
