# deploy/raspberry/systemd/autohunter-scheduler.service
[Unit]
Description=AutoHunter Scheduler (wishlists + notifications)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autohunter
WorkingDirectory=/opt/autohunter
EnvironmentFile=/opt/autohunter/.env
ExecStart=/opt/autohunter/deploy/raspberry/scripts/run_scheduler.sh
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=180
KillMode=mixed

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
